# etc/systemd/system/system-dotfiles-sync.service
[Unit]
Description=Sync system dotfiles

[Service]
Type=oneshot
ExecStart=/home/m0x/.dotfiles/setup-system.sh -r
